# scripts/thepoker.service
[Unit]
Description=Start thepoker GitObs engine

[Service]
Type=oneshot
User=pi
ExecStart={{path-placeholder}}

[Install]
WantedBy=multi-user.target